# package/lean/luci-app-sysmonitor/root/etc/hotplug.d/iface/90-sysmonitor
[ $INTERFACE = wan ] && touch /tmp/sysmonitor
